# deploy/systemd/atlas-primer.service.d/context-reader-graph.conf
# UKG context-reader cutover drop-in for atlas-primer.service.
#
# Install:
#   sudo install -Dm644 \
#     deploy/systemd/atlas-primer.service.d/context-reader-graph.conf \
#     /etc/systemd/system/atlas-primer.service.d/context-reader-graph.conf
#   sudo systemctl daemon-reload
#   sudo systemctl restart atlas-primer
#
# Uninstall (fall back to wiki-glob retrieval):
#   sudo rm /etc/systemd/system/atlas-primer.service.d/context-reader-graph.conf
#   sudo systemctl daemon-reload
#   sudo systemctl restart atlas-primer
#
# When set, the Primer's _context_for_task queries the Unified Knowledge
# Graph (graph.load_latest + is_context_eligible) instead of globbing
# the wiki directory. See docs/UKG_ROLLOUT.md for the full runbook and
# rollback procedure.
[Service]
Environment=EREBUS_CONTEXT_READER=graph
