# infra/local/halyard-keeper/systemd/halyard-keeper.service
# ────────────────────────────────────────────────────────────────
# halyard-keeper.service — user-level systemd unit for the
# Halyard Table Keeper backend.
#
# Install: symlink into ~/.config/systemd/user/, daemon-reload,
# enable --now. See infra/local/halyard-keeper/README.md.
# ────────────────────────────────────────────────────────────────

[Unit]
Description=Halyard Table Keeper backend (halyard-keeper)
Documentation=file:///home/%u/src/agi-hpc/docs/HALYARD_TABLE.md
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
# Pick up NRP_LLM_TOKEN and friends so handle_turn's LLM caller
# has credentials. The "-" prefix means "don't fail if missing"
# — the service gracefully falls back to the keyword stub.
EnvironmentFile=-/home/claude/.primer.env
# Bind loopback only — atlas-caddy fronts the service.
# Tailscale is the operator's OOB admin path, never a service channel.
Environment=HALYARD_KEEPER_HOST=127.0.0.1
Environment=HALYARD_KEEPER_PORT=8091
Environment=LOG_LEVEL=info
# RAG bible path — built from wiki/halyard by build_bible.py.
Environment=HALYARD_BIBLE_PATH=/archive/halyard/bible/halyard_bible.json
# LiveKit: the existing livekit-smoke container runs with --dev
# keys. Override for production.
#
# LIVEKIT_URL is the URL EMBEDDED IN THE MINTED TOKEN RESPONSE —
# i.e. what the *browser* uses to connect. It must be the public
# WSS origin that atlas-caddy fronts, NOT Atlas's loopback. The
# keeper itself doesn't connect to LiveKit; it just signs tokens.
Environment=LIVEKIT_URL=wss://halyard.atlas-sjsu.duckdns.org
Environment=LIVEKIT_API_KEY=devkey
Environment=LIVEKIT_API_SECRET=secret
# Keeper auth: unset in dev (warnings emitted), set in prod.
# Environment=KEEPER_USERNAME=
# Environment=KEEPER_PASSWORD=
# Environment=KEEPER_IP_ALLOWLIST=127.0.0.0/8,100.64.0.0/10
Environment=HALYARD_ARCHIVE_ROOT=/archive/halyard
# handle_turn writes DecisionProof chains to these dirs — redirect
# both AIs' proof chains into /archive/halyard/ so the sandbox's
# ReadWritePaths stays tight.
Environment=ARTEMIS_PROOFS_DIR=/archive/halyard/proofs/artemis
Environment=SIGMA4_PROOFS_DIR=/archive/halyard/proofs/sigma4
Environment=ARTEMIS_SESSIONS_DIR=/archive/halyard/sessions

WorkingDirectory=%h/src/agi-hpc
ExecStart=%h/env/bin/python -m agi.halyard.keeper

Restart=on-failure
RestartSec=5s

# Hardening.
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/archive/halyard
ProtectHome=read-only
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictRealtime=yes

[Install]
WantedBy=default.target
