# deploy/systemd/atlas-artemis-portal.service
[Unit]
Description=Atlas ARTEMIS — Keeper GM portal (HTTP + NATS cache)
Documentation=file:///home/claude/agi-hpc/docs/ARTEMIS.md
After=network-online.target atlas-nats.service
Wants=network-online.target
Requires=atlas-nats.service

[Service]
Type=simple
User=claude
WorkingDirectory=/home/claude/agi-hpc
Environment=PYTHONPATH=/home/claude/agi-hpc/src
Environment=PYTHONUNBUFFERED=1
Environment=ARTEMIS_CHAT_DB=/var/lib/atlas-artemis/chat.sqlite3
Environment=ARTEMIS_PORTAL_PORT=8090
EnvironmentFile=-/home/claude/.artemis.env
ExecStartPre=+/bin/mkdir -p /var/lib/atlas-artemis
ExecStartPre=+/bin/chown claude:claude /var/lib/atlas-artemis
ExecStart=/home/claude/env/bin/python3 -m agi.primer.artemis.portal

Restart=always
RestartSec=10
StartLimitIntervalSec=600
StartLimitBurst=5
RandomizedDelaySec=15

KillSignal=SIGTERM
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal

MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
